# nextcloud-compose.service
[Unit]
Description=Nextcloud
Requisite=docker.service
After=docker.service
Wants=nextcloud-compose-jobs.timer

[Service]
Environment="COMPOSE_DOCKER_CLI_BUILD=1"
ExecStart=/var/lib/nextcloud/nextcloud-compose/nextcloud.sh start
ExecStop=/var/lib/nextcloud/nextcloud-compose/nextcloud.sh stop
RemainAfterExit=true
Restart=on-failure
RestartSec=10
TimeoutSec=infinity
Type=oneshot
WorkingDirectory=/var/lib/nextcloud/nextcloud-compose

[Install]
WantedBy=multi-user.target
