# nextcloud-compose.service
[Unit]
Description=Nextcloud
Requires=docker.service

[Service]
Type=oneshot
RemainAfterExit=true
WorkingDirectory=/opt/nextcloud-compose
ExecStart=/opt/nextcloud-compose/nextcloud.sh start
ExecStop=/opt/nextcloud-compose/nextcloud.sh stop
Restart=on-failure
RestartSec=30

[Install]
WantedBy=multi-user.target
