# nextcloud-compose.service
[Unit]
Description=Nextcloud
After=docker.service

[Service]
WorkingDirectory=/opt/nextcloud-compose
ExecStart=/usr/bin/docker-compose up
ExecStop=/usr/bin/docker-compose down
Restart=on-failure

[Install]
WantedBy=multi-user.target
