# nextcloud-compose.service
[Unit]
Description=Nextcloud
Requisite=docker.service
After=docker.service
Wants=nextcloud-compose-jobs.timer

[Service]
Environment="COMPOSE_DOCKER_CLI_BUILD=1"
ExecStart=/opt/nextcloud-compose/nextcloud.sh start
ExecStop=/opt/nextcloud-compose/nextcloud.sh stop
Restart=on-failure
RestartSec=10
TimeoutSec=infinity
WorkingDirectory=/opt/nextcloud-compose

[Install]
WantedBy=multi-user.target
